# home/powershell/profile.d/path.ps1
# path.ps1 --- Adjust the path to include additional folders
# Mike Barker <mike@thebarkers.com>
# June 30th, 2022

# Get the path variable
Function Get-Pathvariable {
    $env:Path -split ';'
}

# Set the path variable
Function Set-PathVariable {
    param (
        [string]$AddPath,
        [string]$RemovePath
    )
    $regexPaths = @()
    if ($PSBoundParameters.Keys -contains 'AddPath'){
        $regexPaths += [regex]::Escape($AddPath)
    }

    if ($PSBoundParameters.Keys -contains 'RemovePath'){
        $regexPaths += [regex]::Escape($RemovePath)
    }
    
    $arrPath = $env:Path -split ';'
    foreach ($path in $regexPaths) {
        $arrPath = $arrPath | Where-Object {$_ -notMatch "^$path\\?"}
    }
    $env:Path = ($arrPath + $addPath) -join ';'
}
